# recipe.conf
[program:recipe_api]
command=/usr/local/bin/gunicorn -b localhost:8089 -w 4 manage:app
directory=/code
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true